# docker/crontab
# ==========================================================================
# SYNC-HINWEIS
# Diese Datei spiegelt die Run-Types + Zeiten aus .github/workflows/analyze.yml
# und main.py:RUN_TYPES. Bei jeder Aenderung an Cron-Struktur oder Run-Types
# (z.B. Sprint-3-Punkt B: Umbau auf pre_open/post_open, s. PROJECT_STATUS.md)
# MUSS diese Datei mit aktualisiert werden.
#
# Unterschied zu GH Actions: Diese Zeiten sind ECHTE Berlin-Lokalzeit, weil der
# Container TZ=Europe/Berlin gesetzt hat. Kein Sommer/Winter-Drift wie bei
# GH Actions (dort ist der Cron UTC-fix, s. CLAUDE.md Cron-Tabelle) -- lokal
# laeuft es also ganzjaehrig zur "richtigen" Berliner Zeit.
# ==========================================================================
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/local/sbin:/usr/bin:/bin
MAILTO=""

# pre_market -- Mo-Fr 15:00 Berlin
0 15 * * 1-5 root cd /app && python main.py --run-type pre_market >> /app/data/cron.log 2>&1

# evaluate -- Mo-Fr 16:00 Berlin
0 16 * * 1-5 root cd /app && python main.py --run-type evaluate >> /app/data/cron.log 2>&1

# midday -- Mo-Fr 19:00 Berlin
0 19 * * 1-5 root cd /app && python main.py --run-type midday >> /app/data/cron.log 2>&1

# position_check -- Mo-Fr 17:30 Berlin
30 17 * * 1-5 root cd /app && python main.py --run-type position_check >> /app/data/cron.log 2>&1

# close -- Mo-Fr 22:30 Berlin
30 22 * * 1-5 root cd /app && python main.py --run-type close >> /app/data/cron.log 2>&1

# weekly -- So 20:00 Berlin
0 20 * * 0 root cd /app && python main.py --run-type weekly >> /app/data/cron.log 2>&1
